# crontab
# more here
# text here
* * * * * /usr/local/bin/python /fluxObject/fluxer.py /inifiles/AC_chamber.ini ac >> /cronlog 2>&1
# this extra line makes it a valid cron

